# inline.service
[Unit]
Description=Inline Clinic Queue Management System (Debug Logging)
After=graphical.target network-online.target
Wants=graphical.target network-online.target

[Service]
Type=simple
# Delay to ensure the graphical session is ready.
ExecStartPre=/bin/sleep 15
# Change to the correct working directory.
WorkingDirectory=/home/inline/inline
#
# Use sh -c to robustly redirect all stdout and stderr to a log file.
# This is better for debugging than the systemd journal sometimes.
ExecStart=/bin/sh -c "/home/inline/inline/venv/bin/python3 main_program.py >> /home/inline/inline/service.log 2>&1"
#
# Environment variables for the X11 desktop session.
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/inline/.Xauthority"
#
# Restart the service if it fails.
Restart=on-failure
RestartSec=5
# Run as the 'inline' user, who owns the display session.
User=inline

[Install]
WantedBy=graphical.target
